# templates/jetty/jetty.service.erb
#systemd unit file for jetty
#
[Unit]
Description=Jetty Web Application Server
After=syslog.target network.target

[Service]
WorkingDirectory=<%= @idp_jetty_base %>
ExecStart=/usr/bin/java -jar <%= @jetty_home %>/jetty/start.jar jetty.home=<%= @jetty_home %>/jetty jetty.base=<%= @idp_jetty_base %>

User=<%= @jetty_user %>
Group=<%= @jetty_group %>

[Install]
WantedBy=multi-user.target
